# systemd/smart-oomd.service
[Unit]
Description=smart-oomd — kill mémoire préventif par prédiction de tendance
After=default.target

[Service]
Type=simple
WorkingDirectory=/mnt/projects/smart-oomd
ExecStart=/mnt/projects/smart-oomd/.venv/bin/python -m daemon.main
Restart=on-failure
Environment=SMART_OOMD_DRY_RUN=true
Environment=SMART_OOMD_THRESHOLD=15.0
Environment=SMART_OOMD_POLL_INTERVAL=1.0

[Install]
WantedBy=default.target
